# roles/restic/templates/restic-backup.timer
[Unit]
Description=Timer for restic backups

[Timer]
{{ restic_timer_schedule }}

[Install]
WantedBy=timers.target
